# realtek-power-fix.service
[Unit]
Description=Fix Realtek RTL8126 Power Management
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-realtek-power.sh
RemainAfterExit=yes
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
